# robot.rules
SUBSYSTEM=="input",KERNELS=="event5",ACTION=="add",RUN+="/bin/systemctl start robot.service"
